# Pwnable/200-ReadEmail/p200.service.template
[Unit]
Description=Pixels Camp CTF P200 Challenge
After=network.target

[Service]
User=__P200_USER__
Group=__P200_GROUP__
Type=simple
WorkingDirectory=__P200_ROOT__
ExecStart=__P200_ROOT__/p200-runserver.sh
Restart=no

[Install]
WantedBy=multi-user.target
